# shadowsocks.conf
[program:shadowsocks]
command=ssserver -c /etc/shadowsocks/shadowsocks.json
user=root
